# systemd/system/data.automount
[Unit]
Description=Automount /data on WSL

[Automount]
Where=/data
TimeoutIdleSec=0

[Install]
WantedBy=multi-user.target
